# home/dot_config/systemd/user/openrgb-off.service
[Unit]
Description=Turn off all OpenRGB-controlled LEDs at login
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/openrgb --noautoconnect --mode off
RemainAfterExit=yes

[Install]
WantedBy=default.target
